# systemd/tribanft-ipinfo-batch.service
[Unit]
Description=TribanFT IP Geolocation Batch Service
Documentation=https://github.com/n0tjohnny/tribanft
After=network-online.target
Wants=network-online. target

[Service]
Type=simple
User=root
Group=root

# Script path
ExecStart=/usr/bin/python3 /root/tribanft/tools/tribanft-ipinfo-batch.py --daemon

# Auto-restart on failure
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tribanft-ipinfo-batch

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

# Working directory
WorkingDirectory=/root/tribanft
RuntimeDirectory=tribanft

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target